# deployment/systemd/backend.service
[Unit]
Description=AX Stock Backend API
After=network.target

[Service]
Type=simple
Environment=NODE_ENV=production
Environment=APP_TZ=Africa/Kigali
WorkingDirectory=/opt/ax-backend
# Ensure .env exists in WorkingDirectory
ExecStart=/usr/bin/node dist/app.js
Restart=always
RestartSec=5
User=www-data
Group=www-data
# Increase file watchers if needed
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
